# vzlogger/rootfs/etc/services.d/vzlogger/run
#!/usr/bin/with-contenv bashio
# ==============================================================================

CONFIGFILE="$(bashio::config 'configfile')"
bashio::log.info "=============================================================================="
chmod a+x /etc/services.d/vzlogger/listdevices.sh
/etc/services.d/vzlogger/listdevices.sh
bashio::log.info "=============================================================================="
bashio::log.info "$(cat $CONFIGFILE)"
bashio::log.info "=============================================================================="
bashio::log.info "Starting vzlogger"
cp $CONFIGFILE /home/vz/vzlogger.conf
chown vz:vz /home/vz/vzlogger.conf
runuser -l vz -c '/usr/local/bin/vzlogger -c /home/vz/vzlogger.conf --foreground'