# vzlogger/rootfs/etc/services.d/vzlogger/run
#!/usr/bin/with-contenv bashio
# ==============================================================================

CONFIG_PATH=/data/options.json
bashio::log.info "==============================================================================\n"
bashio::log.info "$(cat /data/options.json)"
bashio::log.info "=============================================================================="
bashio::log.info "Starting vzlogger"
runuser -l vz -c '/usr/local/bin/vzlogger -c /data/options.json --foreground'