# vzlogger/rootfs/etc/services.d/vzlogger/run
#!/usr/bin/with-contenv bashio
# ==============================================================================

CONFIGFILE="$(bashio::config 'configfile')"
DEVICE="$(bashio::config 'device')"
bashio::log.info "=============================================================================="
/etc/services.d/vzlogger/listdevices.sh
bashio::log.info "=============================================================================="
bashio::log.info "$(cat $CONFIGFILE)"
bashio::log.info "=============================================================================="
bashio::log.info "Starting vzlogger"
#cp $CONFIGFILE /home/vz/vzlogger.conf
#chown vz:vz /home/vz/vzlogger.conf
#chmod a+w /dev/stdout
#chmod a+rw $DEVICE
#runuser -l vz -c '/usr/local/bin/vzlogger -c /home/vz/vzlogger.conf --foreground'
/usr/local/bin/vzlogger -c $CONFIGFILE --foreground