# vzlogger/rootfs/etc/services.d/vzlogger/run
#!/usr/bin/with-contenv bashio
# ==============================================================================

CONFIG_PATH=/data/options.json
bashio::log.info "=============================================================================="
for sysdevpath in $(find /sys/bus/usb/devices/usb*/ -name dev); do
    (
        syspath="${sysdevpath%/dev}"
        devname="$(udevadm info -q name -p $syspath)"
        [[ "$devname" == "bus/"* ]] && exit
        eval "$(udevadm info -q property --export -p $syspath)"
        [[ -z "$ID_SERIAL" ]] && exit
        bashio::log.info "/dev/$devname - $ID_SERIAL"
    )
done
bashio::log.info "=============================================================================="
bashio::log.info "$(cat /data/options.json)"
bashio::log.info "=============================================================================="
bashio::log.info "Starting vzlogger"
chmod a+r /data/options.json
runuser -l vz -c '/usr/local/bin/vzlogger -c /data/options.json --foreground'