# vzlogger/rootfs/etc/services.d/vzlogger/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Example
# Runs example1 script
# ==============================================================================


bashio::log.info "Starting vzlogger"
CFG=$(bashio::config 'cfg')
bashio::log.info $CFG
echo "$CFG" > /cfg/vzlogger.conf
chmod a+x /cfg/vzlogger.conf
runuser -l vz -c '/usr/local/bin/vzlogger -c /cfg/vzlogger.conf --foreground'